# piano-frontend.service
[Service]
ExecStart=npm run start
WorkingDirectory=%h/piano-appp

[Install]
WantedBy=default.target
