# piano-backend.service
[Service]
ExecStart=python app.py
WorkingDirectory=%h/piano-appp/backend

[Install]
WantedBy=default.target
